# dot_config/systemd/user/runrestic.service
[Unit]
Description=runrestic backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/autorestic backup -a -v
