# dot_config/systemd/user/runrestic.service
[Unit]
Description=runrestic backup

[Service]
Type=oneshot
ExecStart=/usr/bin/python /home/nicolas/.backup_script/backup.py
